# infra/systemd/superagent-agent-service.service
[Unit]
Description=SuperAgent Agent Service
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
ExecStartPre=-/usr/bin/docker rm -f superagent-agent-service
ExecStartPre=-/usr/bin/docker pull eclipse-temurin:21-jre
ExecStart=/usr/bin/docker run --name superagent-agent-service --network host --env-file /opt/superagent/shared/env/agent-service.env -v /opt/superagent/current/agent-service/app.jar:/app/app.jar:ro -v /opt/superagent/current/plugins:/opt/superagent/current/plugins:ro eclipse-temurin:21-jre java -jar /app/app.jar
ExecStop=/usr/bin/docker stop superagent-agent-service
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
